# vpc/snat.service
[Unit]
Description = SNAT via ENI eth0

[Service]
ExecStart = /opt/nat/snat.sh
Type = oneshot

[Install]
WantedBy = multi-user.target